# cloud-bot/com.tg-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tg-bot</string>
    <key>ProgramArguments</key>
    <array>
        <!-- TODO: Update these paths for your machine -->
        <string>/path/to/.bun/bin/bun</string>
        <string>run</string>
        <string>/path/to/cloud-bot/bot.ts</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/path/to/home</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/path/to/home</string>
        <key>PATH</key>
        <string>/path/to/.local/bin:/path/to/.bun/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>/path/to/logs/tg-bot.log</string>
    <key>StandardErrorPath</key>
    <string>/path/to/logs/tg-bot.log</string>
    <key>AbandonProcessGroup</key>
    <true/>
</dict>
</plist>
